# ups-agent.service
[Unit]
Description=UPS Cluster Agent
After=network-online.target nut-server.service
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/bin/bash /opt/ups-agent/ups_agent.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
